# roles/factorio-server/files/supervisor.conf
[program:factorio]
user=factorio
command=/opt/factorio/start.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/factorio.err.log
stdout_logfile=/var/log/factorio.out.log
stopsignal=QUIT
